# overlay/opt/zhiverbox/scripts/etc/systemd/system/tor-motd-info.timer
[Unit]
Description=Starts tor-motd-info.service after some delay
After=tor-motd-info.path

[Timer]
OnActiveSec=20
Unit=tor-motd-info.service
