# overlay/etc/systemd/system/sshd.service.d/override.conf
[Service]
# set affinity of sshd to the LITTLE cores (0-3)
ExecStart=/usr/bin/taskset -a -c 0-3 /usr/sbin/sshd -D $SSHD_OPTS
